# root/etc/uci-defaults/luci-app-fleth
#!/bin/sh

chmod +x /usr/sbin/fleth /etc/hotplug.d/iface/25-fleth-ipprefix

exit 0
