# root/etc/uci-defaults/luci-app-fleth
#!/bin/sh

chmod +x /usr/sbin/fleth /etc/init.d/fleth /etc/hotplug.d/iface/25-fleth-ip6prefix

if [ ! -f /etc/config/fleth ]; then
    cat <<EOF > /etc/config/fleth
config fleth 'global'
	option enabled '0'
    option ip6prefix_enabled '0'
	option type 'auto'
	option interface 'wan'
	option interface6 'wan6'
	option mtu '1460'

EOF
fi

exit 0
