# root/etc/uci-defaults/luci-app-fleth
#!/bin/sh

chmod +x /usr/sbin/fleth /etc/init.d/fleth /etc/hotplug.d/iface/70-fleth

if [ ! -f /etc/config/fleth ]; then
    cat <<EOF > /etc/config/fleth
config fleth 'global'
	option enabled '0'
	option interface 'wan'
	option interface_zone 'wan'
	option interface6 'wan6'
	option mtu '1460'
	option prefer_slaac '1'
	option tunnel_activation '1'
	option upnp_auto_config '1'
EOF
else
    sed -i '/option ip6prefix_enabled/d' /etc/config/fleth
    sed -i '/option cron_dhcpv6_renew_enabled/d' /etc/config/fleth
    sed -i '/option type/d' /etc/config/fleth
    # Add prefer_slaac option if not exists
    if ! grep -q "option prefer_slaac" /etc/config/fleth; then
        uci set fleth.global.prefer_slaac='1'
        uci commit fleth
    fi
    # Add tunnel_activation option if not exists
    if ! grep -q "option tunnel_activation" /etc/config/fleth; then
        uci set fleth.global.tunnel_activation='1'
        uci commit fleth
    fi
    # Add upnp_auto_config option if not exists
    if ! grep -q "option upnp_auto_config" /etc/config/fleth; then
        uci set fleth.global.upnp_auto_config='1'
        uci commit fleth
    fi
fi
fleth hook_none.js
exit 0
